# omniclone.timer
[Unit]
Description=Run Rclone Omniscript Service every 30 minutes

[Timer]
OnBootSec=5min
OnUnitInactiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
